# deploy/everyonemessage-api.service
[Unit]
Description=EveryoneMessage API
After=network.target postgresql.service redis-server.service

[Service]
Type=simple
User=emapp
Group=emapp
WorkingDirectory=/opt/everyonemessage
Environment=NODE_ENV=production
ExecStart=/usr/bin/node apps/api/dist/index.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
